# demo_swarm/tests/etc/supervisord-controlpanel.conf
[supervisord]
logfile = /data/logs/supervisord.log
logfile_maxbytes = 50MB
logfile_backups=10
loglevel = info
pidfile = /data/pids/supervisord.pid
nodaemon = false
minfds = 1024
minprocs = 200
umask = 022
identifier = supervisor
directory = /tmp
nocleanup = true
childlogdir = /tmp
strip_ansi = false

[program:node]
directory = /data
nodaemon = true
command=/usr/bin/node controlpanel/app.js
logfile = /data/logs/node.log
logfile_maxbytes = 50MB
logfile_backups=0
loglevel = info
